# openstack/mistral/mistral-executor.service
[Unit]
Description=Mistral Executor Server
After=syslog.target network.target

[Service]
Type=simple
User=mistral
Restart=on-failure
ExecStart=/usr/bin/mistral-server --config-file=/etc/mistral/mistral.conf --server=executor

[Install]
WantedBy=multi-user.target
