# openstack/swift/openstack-swift-object.service
[Unit]
Description=OpenStack Swift Object Server
After=syslog.target network.target

[Service]
Type=simple
User=swift
Group=swift
WorkingDirectory=/var/lib/swift
PrivateTmp=yes
ExecStart=/usr/bin/swift-object-server /etc/swift/object-server.conf

[Install]
WantedBy=multi-user.target
