# openstack/cinder/openstack-cinder-scheduler.service
[Unit]
Description=OpenStack Cinder Scheduler Server
After=syslog.target network.target

[Service]
Type=notify
User=cinder
Group=cinder
WorkingDirectory=/var/lib/cinder
PrivateTmp=yes
ExecStart=/usr/bin/cinder-scheduler
Restart=on-failure

[Install]
WantedBy=multi-user.target
